# Daily-Brief/morning-brief/deploy/systemd/morning-brief.timer.example
[Unit]
Description=Run daily morning brief at 7:30 AM Eastern

[Timer]
# Keep the delivery time tied to Eastern time even if the server runs in UTC.
OnCalendar=*-*-* 07:30:00 America/New_York
Persistent=true

[Install]
WantedBy=timers.target
